# upstart/host/gv-valley-host-runner.conf
start on starting gv-valley-host
stop on stopping gv-valley-host
