# upstart/host/gv-valley-host-runner-1.conf
start on starting gv-valley-host-runner
stop on stopping gv-valley-host-runner
respawn

exec su - root -c 'gv-valley provide gv/valley/runner >> /var/log/gv-valley-host/runner-1.log 2>&1'
