# upstart/bedrock/gv-bedrock-config.conf
start on started gv-bedrock-server
stop on stopping gv-bedrock