# upstart/bedrock/gv-bedrock-server-1.conf
start on starting gv-bedrock-server
stop on stopping gv-bedrock-server
respawn

exec su - root -c 'gv-bedrock provide gv/bedrock/server >> /var/log/gv-bedrock/server-1.log 2>&1'
