# upstart/bedrock/gv-bedrock-config-1.conf
start on starting gv-bedrock-config
stop on stopping gv-bedrock-config
respawn

exec su - root -c 'gv-bedrock provide gv/bedrock/config >> /var/log/gv-bedrock/config-1.log 2>&1'
