# upstart/bedrock/gv-bedrock-server.conf
start on starting gv-bedrock
stop on stopping gv-bedrock
