# upstart/addon/gv-valley-addon.conf
pre-start script

bash << "EOF"
  mkdir -p /var/log/gv-valley-addon
  chown -R root /var/log/gv-valley-addon
EOF

end script

start on runlevel [2345]

stop on runlevel [016]
