# upstart/addon/gv-valley-addon-a_etcd-1.conf
start on starting gv-valley-addon-a_etcd
stop on stopping gv-valley-addon-a_etcd
respawn

exec su - root -c 'gv-valley provide gv/addons/etcd >> /var/log/gv-valley-addon/a_etcd-1.log 2>&1'
