# upstart/addon/gv-valley-addon-a_postgresql.conf
start on starting gv-valley-addon
stop on stopping gv-valley-addon
