# upstart/valley/gv-valley-file_server-1.conf
start on starting gv-valley-file_server
stop on stopping gv-valley-file_server
respawn

exec su - root -c 'gv-valley provide gv/valley/file_server >> /var/log/gv-valley/file_server-1.log 2>&1'
