# upstart/valley/gv-valley-etcd-1.conf
start on starting gv-valley-etcd
stop on stopping gv-valley-etcd
respawn

exec su - root -c 'gv-valley provide gv/valley/etcd >> /var/log/gv-valley/etcd-1.log 2>&1'
