# upstart/valley/gv-valley-deployer-1.conf
start on starting gv-valley-deployer
stop on stopping gv-valley-deployer
respawn

exec su - root -c 'gv-valley provide gv/valley/deployer >> /var/log/gv-valley/deployer-1.log 2>&1'
