# upstart/valley/gv-valley-file_server.conf
start on starting gv-valley
stop on stopping gv-valley
