# upstart/valley/gv-valley-balancer-1.conf
start on starting gv-valley-balancer
stop on stopping gv-valley-balancer
respawn

exec su - root -c 'gv-valley provide gv/valley/balancer >> /var/log/gv-valley/balancer-1.log 2>&1'
